# daisyproducer/dictionary/misc/cron.daily/daisyproducerWriteLiblouisTables
#!/bin/sh

INSTALLDIR=/path/to/daisyproducer
umask 113

# Write the liblouis tables
cd $INSTALLDIR
python manage.py writeLiblouisTables
