# dictionary/misc/cron.hourly
#!/bin/sh

INSTALLDIR=/path/to/daisyproducer

# Write the liblouis tables
cd $INSTALLDIR
python manage.py writeLiblouisTables